# src/paloma/templates/paloma/conf/etc/supervisord.conf.d/conf.ini
[program:sluworker]
command={{ VIRTUAL_ENV}}/bin/celery_workman.py {{ PROJECT_DIR }} start
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
user=hdknr
directory={{ PROJECT_DIR }}
